# transcript-processing/com.exobrain.supernote-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.exobrain.supernote-watcher</string>

    <key>RunAtLoad</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alexhedtke/Documents/Exobrain harness/transcript-processing/run-process-supernote.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/alexhedtke/My Drive/Supernote/Note</string>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/exobrain-supernote-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/exobrain-supernote-watcher.err</string>

    <!-- Periodic fallback in case WatchPaths misses a trigger
         (e.g., Google Drive mounts after launchd already checked) -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
