# discord/com.exobrain.discord-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.exobrain.discord-digest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alexhedtke/Documents/Exobrain harness/discord/run-discord-digest.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/alexhedtke/Documents/Exobrain harness/discord</string>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/alexhedtke/.claude/channels/discord/digest-fetch.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alexhedtke/.claude/channels/discord/digest-fetch-error.log</string>
</dict>
</plist>
