# apple-notes-sync/com.exobrain.apple-notes-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.exobrain.apple-notes-sync</string>

    <!-- Source of truth: /Users/alexhedtke/Documents/Exobrain harness/apple-notes-sync/
         Deployed to ~/apple_notes_sync/ because launchd can't read ~/Documents/ (TCC) -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/alexhedtke/apple_notes_sync/apple_notes_sync.py</string>
        <string>--config</string>
        <string>/Users/alexhedtke/apple_notes_sync/apple_notes_sync_config.json</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/apple_notes_sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/apple_notes_sync_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
