# things3-sync/com.exobrain.things3-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.exobrain.things3-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alexhedtke/Documents/Exobrain harness/things3-sync/run-things3-sync.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/tmp/exobrain-things3-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/exobrain-things3-sync.err</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
